# crates/powerutil/powerutil.timer
[Unit]
Description=Run PowerUtil every few seconds

[Timer]
OnUnitActiveSec=5s
Unit=powerutil.service

[Install]
WantedBy=timers.target