# crates/powerutil/powerutil.timer
[Unit]
Description=Run PowerUtil every few mins

[Timer]
OnUnitActiveSec=300s
Unit=powerutil.service

[Install]
WantedBy=timers.target